# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/marketphase/bin/console schedule:run >> /dev/null 2>&1
